# cmd/appherder/systemd/appherder-sync.service
[Unit]
Description=Reconcile AppImages with installed applications

[Service]
Type=oneshot
ExecStart={{BIN}} sync
