# cmd/appherder/systemd/appherder-sync.path
[Unit]
Description=Watch AppImages for appherder sync

[Path]
PathChanged={{APPIMAGES_DIR}}
Unit=appherder-sync.service

[Install]
WantedBy=default.target
